# dockermail/clamav/supervisord.d/clamav.ini
[program:clamav]
process_name = clamav
command=clamd
startsecs=0
autorestart=true

[program:freshclam]
process_name = freshclam
command=freshclam
startsecs=0
autorestart=false

[program:clamav_milter]
process_name = clamav_milter
command=clamav-milter
startsecs=0
autorestart=true
