# dockermail/postfix/supervisord.d/postfix.ini
[program:postfix]
process_name = postfix_startup
command=postfix start
startsecs=0
autorestart=false

[program:opendkim]
process_name = opendkim_startup
command=opendkim
startsecs=0
autorestart=false

[program:opendmarc]
process_name = opendmarc_startup
command=opendmarc -p unix:/var/run/opendmarc/opendmarc.sock
startsecs=0
autorestart=false

[program:opendkim_perms]
process_name = opendkim_perms
command=chown opendkim:opendkim /etc/opendkim/mail
startsecs=0
autorestart=false

[program:opendmarc_perms]
process_name=opendmarc_perms
command=chown postfix:postfix /run/opendmarc/opendmarc.sock
startsecs=0
autorestart=false
