# dockermail/postfix/supervisord.d/postfix.ini
[program:postfix]
process_name = postfix_startup
command=postfix start
startsecs=0
autorestart=false

[program:opendkim]
process_name = opendkim_startup
command=opendkim
user=opendkim
startsecs=0
autorestart=false

[program:opendmarc]
process_name = opendmarc_startup
command=opendmarc -p unix:/var/run/opendmarc/opendmarc.sock
user=postfix
startsecs=0
autorestart=false

[program:newaliases]
process_name=newaliases
command=newaliases
startsecs=0
autorestart=false
